# roles/borg/templates/borg-backup.timer.j2
[Unit]
Description=Run Borg Backup Hourly
Requires=borg-backup.service

[Timer]
OnCalendar=hourly
Persistent=true

[Install]
WantedBy=timers.target